# logsnarf-rb/ops/templates/restart-logsnarf.service
[Unit]
Description=Restart logsnarf service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl try-restart logsnarf.service
