# config/logsnarf.socket
[Unit]
Description=Logsnarf listen socket

[Socket]
ListenStream=443
NoDelay=true
